# kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
