# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
